# cron/crontab
1 0 * * * /home/ec2-user/youtube-export/cron/convert.sh
1 12 * * * /home/ec2-user/youtube-export/cron/publish.sh
1 6 */2 * * /home/ec2-user/youtube-export/cron/language_data.sh 
