# cron/crontab
* 0 * * * /home/ec2-user/youtube-export/cron/convert.sh
* 12 * * * /home/ec2-user/youtube-export/cron/publish.sh
